# telegram-carousel-bot.service
[Unit]
Description=Telegram Carousel Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/3percents
EnvironmentFile=/etc/telegram-carousel-bot.env
ExecStart=/opt/3percents/.venv/bin/python -m src.bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
